# pkg-plist
/usr/local/opnsense/scripts/smt/smt_plugin.py
/usr/local/opnsense/scripts/smt/requirements.txt
/usr/local/opnsense/scripts/smt/config.yaml
/var/log/smt_plugin.log
/etc/rc.d/smt
/usr/local/opnsense/mvc/app/controllers/OPNSense/SMT/ConfigController.php
/usr/local/opnsense/mvc/app/models/OPNSense/SMT/Config.php
/usr/local/opnsense/mvc/app/templates/OPNSense/SMT/Config/config.ui.php
/usr/local/opnsense/mvc/app/javascript/smt.js
/usr/local/share/examples/smt/social_media_domains.txt
/usr/local/share/examples/smt/config.yaml

